# scheduled_login_tradetron.service
[Unit]
Description=Login into tradetron via selenium

[Service]
ExecStart=python3 -m /home/ubuntu/bin/scheduled_login_tradetron/app
Type=oneshot
