# systemd/notify_failed@scheduled_login_tradetron.service
[Unit]
Description=status email for %I to user

[Service]
Type=oneshot
ExecStart=cat /home/ubuntu/apps/scheduled_login_tradetron/logs/scheduled_login_tradetron.log | mail -s "notify_failed@scheduled_login_tradetron" sriv.ritesh@gmail.com
